# raspberry-pi/deploy/eval-teacher.timer
# Fires eval-teacher.service nightly. Persistent=true runs a missed
# window at next boot, and systemd never starts a second instance while
# one is still running.

[Unit]
Description=Nightly eval teacher batch

[Timer]
OnCalendar=*-*-* 03:00
Persistent=true
# Spread start over a few minutes so a future second Pi doesn't sync-hit
# the cloud API at the same second.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
