# raspberry-pi/deploy/receipt-sync.timer
# Fires receipt-sync.service every minute — cheap, because the sync makes
# no cloud calls when nothing new is pending. Persistent=true runs a
# missed window at next boot, and systemd never starts a second instance
# while one is still running.

[Unit]
Description=Periodic receipt log cloud sync

[Timer]
OnCalendar=minutely
Persistent=true
# Small jitter so a future second Pi doesn't sync-hit the cloud API at
# the same second; kept well under the interval.
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
